# vimwiki/Library/LaunchAgents/phajas.vimwiki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>phajas.vimwiki</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/phajas/bin/vimwiki_update</string>
	</array>
	<key>StartInterval</key>
	<integer>60</integer>
</dict>
</plist>
